# com.docket.social.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.docket.social.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/gregorybennett/docket-social/.venv/bin/python -m src.main scrape &amp;&amp; /Users/gregorybennett/docket-social/.venv/bin/python -m src.main schedule --auto</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gregorybennett/docket-social</string>

    <!-- Monday at 12:00 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>12</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gregorybennett/docket-social/logs/launchd.weekly.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gregorybennett/docket-social/logs/launchd.weekly.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
