# com.docket.social.token-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.docket.social.token-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gregorybennett/docket-social/.venv/bin/python</string>
        <string>-m</string>
        <string>src.main</string>
        <string>refresh-tokens</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gregorybennett/docket-social</string>

    <!-- 1st and 15th of each month at 10:00 AM (twice-monthly for safety) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Day</key><integer>1</integer>
            <key>Hour</key><integer>10</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Day</key><integer>15</integer>
            <key>Hour</key><integer>10</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/gregorybennett/docket-social/logs/launchd.token-refresh.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/gregorybennett/docket-social/logs/launchd.token-refresh.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
